# local_llm_router/systemd/local-llm-router.service.example
[Unit]
Description=local_llm_router daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/local-llm-router/app
ExecStart=%h/.local/bin/local-llm-router serve
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
